# meta-ppe-app/recipes-ppe/ppe-infer/files/ppe-detect-static.service
[Unit]
Description=PPE Static ONNX Detection (one-shot test) 
After=network.target

[Service]
Type=simple
# ví dụ chạy test trên 1 ảnh mẫu có sẵn trong /opt/models/static_img.jpg
ExecStart=/usr/bin/detect_static_onnx.py \
  --model /opt/models/ppe_model.onnx \
  --image /opt/models/static_img.jpg \
  --out   /var/volatile/lib/ppe/annotated.jpg \
  --log   /var/volatile/log/ppe/detections.jsonl \
  --classes /opt/models/classes.txt \
  --cpu

# Nếu muốn lặp lại định kỳ 30s (demo), dùng:
# ExecStart=/bin/sh -c 'while true; do /usr/bin/detect_static_onnx.py --model /opt/models/ppe_model.onnx --image /opt/models/static_img.jpg --out /var/lib/ppe/annotated.jpg --log /var/log/ppe/detections.jsonl --classes /opt/models/classes.txt --cpu; sleep 30; done'

Restart=on-failure
RestartSec=2
User=root

[Install]
WantedBy=multi-user.target
